# deploy/demo/systemd/lattice-demo-reset.service
# Installed by demo-bootstrap.sh with __REPO__ substituted. One reset run;
# scheduling lives in lattice-demo-reset.timer.
[Unit]
Description=Lattice hosted demo — nightly world reset (wipe + reseed)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=__REPO__
ExecStart=__REPO__/deploy/demo/demo-reset.sh
# Covers the wipe, stack rebuild, package installs and showcase seed. The
# post-wipe auth-plane rescan itself (tens of minutes to multiple hours on a
# 2 vCPU box) runs in the background and outlives this unit; demo-up.sh's own
# demo-Loupe attempt is bounded (PROVISION_TIMEOUT_SECONDS=90) rather than
# waiting it out, and lattice-demo-loupe-retry.timer picks the Loupe start up
# later, whenever the rescan actually drains (F21 Fire 3).
TimeoutStartSec=1800
# demo-up.sh leaves the stack's host processes (processor … facet) running in
# the background, lifecycle-managed by the scripts' own pkills; the default
# control-group cleanup would reap them all the moment the oneshot deactivates.
KillMode=process
